# configs/alsa-restore-delayed.timer
[Unit]
Description=Delayed ALSA mixer restore after boot

[Timer]
OnBootSec=10s

[Install]
WantedBy=timers.target